# launchd/com.timelapse.healthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Hammerspoon healthcheck (5分間隔)
  heartbeat ファイルが stale なら hs.reload() を実行
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.timelapse.healthcheck</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jyajyamaru/.timelapse/timelapse-healthcheck.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
